# HomeBrain.cron
#!/bin/bash

cd "$(dirname "$0")";
DIR="$( cd "$( dirname "${BASH_SOURCE[0]}" )" && pwd )"

$DIR/_tools/dbBackup.php

$DIR/lan/dbBackup.php

$DIR/heating/getTempSet.php
$DIR/heating/dbBackup.php


/sbin/shutdown -r now
